# apps/worker/bidradar-worker.service
[Unit]
Description=Bidradar Worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bidradar
WorkingDirectory=/opt/bidradar
ExecStart=/usr/bin/node /opt/bidradar/worker/dist/index.js
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
